# frps.service
[Unit]
Description=frps
After=multi-user.target
 
[Service]
TimeoutStartSec=30
ExecStart=/etc/frp/frps -c /etc/frp/frps.ini
ExecStop=/bin/kill $MAINPID
 
[Install]
WantedBy=multi-user.target
